# scheduler/org.cups.cupsd.service.in
[Unit]
Description=CUPS Scheduler
Documentation=man:cupsd(8)
After=sssd.service

[Service]
ExecStart=@sbindir@/cupsd -l
Type=simple

[Install]
Also=cups.socket cups.path
WantedBy=printer.target
